# launchd/com.girl-agent.sillytavern.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.girl-agent.sillytavern</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/geoff/Projects/Girl-Agent/scripts/run_sillytavern.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/geoff/Projects/Girl-Agent</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/geoff/Projects/Girl-Agent/logs/sillytavern.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/geoff/Projects/Girl-Agent/logs/sillytavern.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/geoff/homebrew/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
